# config/systemd/user/akeyshually.service
[Unit]
Description=akeyshually - Minimal keyboard shortcut daemon
After=graphical-session.target kanata.service
PartOf=kanata.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=%h/.local/bin/akeyshually
Restart=always
RestartSec=3
KillMode=process

[Install]
WantedBy=default.target
